# scripts/rc.local
#!/bin/sh -e
docker run -d -p 8050:8050 -p 5023:5023 scrapinghub/splash --disable-ui --disable-lua
